# server/deploy/threadkit-http.service
[Unit]
Description=ThreadKit HTTP Server
After=network.target

[Service]
Type=simple
User=threadkit
Group=threadkit
WorkingDirectory=/opt/threadkit
EnvironmentFile=/opt/threadkit/.env
ExecStart=/opt/threadkit/bin/threadkit-http
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=threadkit-http

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/threadkit

[Install]
WantedBy=multi-user.target
